# config/files/usr/etc/systemd/system/deeptank.mount
[Mount]
What=//192.168.10.2/creed
Where=/home/creed/deeptank
Options=credentials=/home/creed/deeptank.credentials,noperm,vers=4,_netdev
Type=cifs
